# config/docker/cron/crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/30 * * * * root cd /var/www/html && bin/drush cron >/dev/null 2>&1

